# intel/ops/market-intel-watchdog.timer
[Unit]
Description=Run the market-intel stale-heartbeat watchdog every 10 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
